# defaults/.config/polybar/scripts/95-usb.rules
KERNEL=="sd*", ACTION=="add", ATTR{removable}=="1", \
    RUN+="/home/reinaldo/.config/polybar/scripts/system-usb-udev.sh --update"
KERNEL=="sd*", ACTION=="remove", \
    RUN+="/home/reinaldo/.config/polybar/scripts/system-usb-udev.sh --update"

