# defaults/.config/systemd/user/clock-in-reminder.timer
[Unit]
Description=Clock-in Reminder Timer

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target

