# defaults/.config/systemd/user/mbsync.service
[Unit]
Description=Mailbox synchronization service

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync -Vac /home/reinaldo/.config/isync/mbsyncrc
