# defaults/.config/systemd/user/wallpaper.timer
[Unit]
Description=Change wallpapers every 30 minutes

[Timer]
OnActiveSec=5sec
OnUnitActiveSec=24h

[Install]
WantedBy=timers.target
